# extra/esajpip.service
[Unit]
Description=Enables esajpip

[Service]
ExecStart=/usr/bin/esajpip
User=user
Group=group

[Install]
WantedBy=multi-user.target